# tool/cron/updateAlgoData.cron
* */2 * * * /opt/wpp/src/wlan-pos/tool/cron/updateAlgoData.sh
